# cloud/services/nexus-leadgen-api.service
[Unit]
Description=NEXUS Lead Gen API (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__NEXUS_USER__
Group=__NEXUS_USER__
WorkingDirectory=__PROJECT_DIR__/lead-gen
ExecStart=__VENV_DIR__/bin/uvicorn src.main:app --host 127.0.0.1 --port 8082 --log-level warning
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=__ENV_FILE__

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexus-leadgen-api

# Security hardening
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=__PROJECT_DIR__

[Install]
WantedBy=multi-user.target
